# server-bots/ybseo-content-bot/ybseo-content-bot.timer
[Unit]
Description=Check whether YBSEO needs a new article

[Timer]
OnCalendar=*-*-* 08:15:00 Europe/Istanbul
Persistent=true
RandomizedDelaySec=30m

[Install]
WantedBy=timers.target
